# calibre-web-automated/rootfs/etc/s6-overlay/s6-rc.d/init-tailscale-up/up
/etc/s6-overlay/scripts/init-tailscale-up
